# pardus-greeter-ebaqr/DEBIAN/postinst
#!/bin/sh
systemctl enable pardus-greeter-ebaqr.service
chmod 755 /etc/qrpsw


