# debian/postinst
#!/bin/sh
systemctl enable pardus-lightdm-greeter-listener.service||true
chmod 755 /etc/qrpsw||true


