# pardus-greeter-open/DEBIAN/postinst
#!/bin/sh
systemctl enable pardus-greeter-open.service
chmod 755 /etc/qrpsw


