# pardus-greeter-open/DEBIAN/postinst
#!/bin/sh
systemctl enable pardus-greeter-open.service


